# roles/managed_node/templates/snapraid/snapraid-runner.timer.j2
[Unit]
Description=Snapraid runner

[Timer]
RandomizedDelaySec=10
OnCalendar=*-*-* 20:00:00

[Install]
WantedBy=timers.target
