# roles/managed_node/templates/snapraid/snapraid-runner.timer.j2
[Unit]
Description=Snapraid runner

[Timer]
RandomizedDelaySec=10
OnCalendar={{ snapraid_runner_timer }}

[Install]
WantedBy=timers.target