# roles/managed_node/templates/system/shutdown.timer.j2
[Unit]
Description={{ template_var.systemd_description }}

[Timer]
RandomizedDelaySec=10
OnCalendar={{ template_var.timer }}

[Install]
WantedBy=timers.target