# roles/managed_node/templates/system/shutdown.service.j2
[Unit]
Description={{ template_var.systemd_description }}
Wants=network.target
After=network-online.target

[Service]
ExecStart=shutdown
Type=oneshot

[Install]
WantedBy=multi-user.target
