# roles/managed_node/templates/container/ansible.service.j2
[Unit]
Description={{ systemd_description }}
Wants=network.target
After=network-online.target
RequiresMountsFor=%h/.local/share/containers/storage %t/containers

[Service]
WorkingDirectory={{ systemd_working_directory }}
ExecStart={{ command }}
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=default.target