# roles/managed_node/templates/container/ansible.timer.j2
[Unit]
Description={{ timer_description }}

[Timer]
RandomizedDelaySec=10
OnCalendar={{ timer }}
Persistent=true

[Install]
WantedBy=timers.target