# roles/managed_node/templates/user/template.timer.j2
[Unit]
Description={{ template_var.systemd_description }}

[Timer]
RandomizedDelaySec=10
OnCalendar={{ template_var.timer }}
Persistent=true

[Install]
WantedBy=timers.target
