# roles/control_node/files/container_systemd/image-prune/container-image-prune_sample.service
[Unit]
Description=Prune container images
Wants=network.target
After=network-online.target
RequiresMountsFor=%h/.local/share/containers/storage %t/containers

[Service]
ExecStart=podman image prune --force
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=default.target