# roles/control_node/files/container_systemd/update/container-update_sample.timer
[Unit]
Description=Run update by time frame

[Timer]
RandomizedDelaySec=10
OnCalendar=CalendarTimer
Persistent=true

[Install]
WantedBy=timers.target