# roles/managed_node/templates/restic/restic-unlock.service.j2
[Unit]
Description=Restic backup service
[Service]
Type=oneshot
ExecStart=restic unlock
EnvironmentFile={{ restic_env_path }}
