# roles/container_systemd/files/container-update_sample.timer
[Unit]
Description=Run update by time frame

[Timer]
OnCalendar=CalendarTimer
Persistent=true

[Install]
WantedBy=timers.target