# roles/control_node/templates/container_systemd/image-prune/container-image-prune.timer.j2
[Unit]
Description=Clean container images

[Timer]
RandomizedDelaySec=10
OnCalendar=CalendarTimer
Persistent=true

[Install]
WantedBy=timers.target