# roles/managed_node/templates/systemd/shutdown.service.j2
[Unit]
Description=Shutdown system
Wants=network.target
After=network-online.target

[Service]
ExecStart=shutdown
Type=oneshot

[Install]
WantedBy=multi-user.target
