# roles/systemd_timer/templates/template.timer.j2
[Unit]
Description={{ systemd_timer_name }}

[Timer]
RandomizedDelaySec=10
OnCalendar={{ systemd_timer_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
