# roles/systemd_timer/templates/template.timer.j2
[Unit]
Description={{ systemd_description }}

[Timer]
RandomizedDelaySec=10
OnCalendar={{ systemd_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
