# roles/startup/files/container-startup_sample.service
[Unit]
Description=Run container service at startup
Wants=network.target
After=network-online.target
RequiresMountsFor=$HOME/.local/share/containers/storage $XDG_RUNTIME_DIR/containers

[Service]
WorkingDirectory=ServiceDirectory
ExecStart=/usr/bin/ansible-playbook ServiceDirectory/container.yaml
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=default.target